# outflits-erp/backend/deploy/supervisor-outflits-worker.conf
[program:outflits-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/outflits-erp/outflits-erp/backend/artisan queue:work redis --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/supervisor/outflits-worker.log
stopwaitsecs=3600
